# skills/kernel-observer/openclaw-kernel-observer.service
[Unit]
Description=OpenClaw Kernel Observer Bridge
Documentation=https://docs.openclaw.ai/skills/kernel-observer
After=network.target openclaw-gateway.service
Requires=openclaw-gateway.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/csilva/.openclaw/workspace/skills/kernel-observer
ExecStart=/usr/bin/node /home/csilva/.openclaw/workspace/skills/kernel-observer/bridge.js --scope=all
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_BPF CAP_PERFMON
AmbientCapabilities=CAP_SYS_ADMIN CAP_BPF CAP_PERFMON
NoNewPrivileges=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target